# stage75/00-install-greengrass/files/greengrasscertcopy.service
[Unit]
Description=Greengrass Certificate Copy at Boot

[Service]
ExecStart=/greengrass/greengrasscertcopy.sh

[Install]
WantedBy=multi-user.target
