# stage60-wifi-supplicant/00-install-wpa-supplicant/files/wifisupplicantcopy.service
[Unit]
Description=Copy WiFi Supplicant at Boot

[Service]
ExecStart=/opt/wifisupplicantcopy.sh

[Install]
WantedBy=multi-user.target
